# diverses/crontab
0 0 * * * /home/pi/Temperaturmessung/diverses/mitternacht.sh >/dev/null 2>&1
*/10 * * * * /home/pi/Temperaturmessung/diverses/gnuplot.sh
*/10 * * * * /home/pi/Temperaturmessung/diverses/ftp.sh
