# diverses/crontab
0 0 * * * /home/pi/Temperaturmessung/diverses/mitternacht.sh
